# debian/corecluster.cron.d
* *     * * *	root	[ -x /usr/sbin/cc-manage ] && /usr/sbin/cc-manage cron minute
31 *     * * *	root	[ -x /usr/sbin/cc-manage ] && /usr/sbin/cc-manage cron hourly
1 0    * * *	root	[ -x /usr/sbin/cc-manage ] && /usr/sbin/cc-manage cron daily
31 3    * * 7	root	[ -x /usr/sbin/cc-manage ] && /usr/sbin/cc-manage cron weekly
22 4    1 * *	root	[ -x /usr/sbin/cc-manage ] && /usr/sbin/cc-manage cron monthly
